# deploy/systemd/fileflow-web.service
[Unit]
Description=FileFlow Flask Web (Gunicorn)
After=network.target fileflow-socket.service
Requires=fileflow-socket.service

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_DEBUG=0
ExecStart=__APP_DIR__/.venv/bin/gunicorn --workers 2 --threads 4 --bind 127.0.0.1:5000 wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target